# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
loglevel=Critical
logfile_maxbytes=0
user=root

[program:tomcat]
command=/usr/share/tomcat9/bin/catalina.sh run
autostart=true
autorestart=true

